# #ubuntu.dyn-o-saur.com/@Mac_OSX/Library/LaunchDaemons/monit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" ">
<plist version="1.0">
<dict>
        <key>Disabled</key>
        <false/>
        <key>Label</key>
        <string>com.vyperlogix.monit</string>
        <key>OnDemand</key>
        <false/>
        <key>ProgramArguments</key>
        <array>
                <string>/opt/monit-5.3/bin/launchd_wrapper.sh</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>ServiceDescription</key>
        <string>Monit</string>
        <key>StandardErrorPath</key>
        <string>/opt/monit-5.3/logs/launchd.stderr</string>
        <key>StandardOutPath</key>
        <string>/opt/monit-5.3/logs/launchd.stdout</string>
        <key>UserName</key>
        <string>root</string>
</dict>
</plist>